# make_ams_client.service
[Unit]
Description=MAKE; AMS Client
Documentation=https://github.com/szkotgh/make-ams_client.git
After=network.target

[Service]
User=administrator
Environment=DISPLAY=:0
ExecStart=/bin/bash /home/administrator/make-ams_client/stater.sh
Restart=on-failure
RestartSec=1

[Install]
WantedBy=graphical.target
